# systemd/hisnos-threat-engine.service
[Unit]
Description=HisnOS Threat Detection Engine
Documentation=https://hisnos.dev/docs/threat-engine
# HARDENED: threat engine must not block boot or user session.
After=graphical.target hisnosd.service auditd.service nftables.service
Requires=hisnosd.service
BindsTo=hisnosd.service
ConditionPathExists=/etc/hisnos/release

[Service]
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
NoNewPrivileges=yes
# HARDENED: Type=idle — only start when nothing else needs CPU.
Type=idle
ExecStart=/usr/local/bin/hisnosd --module threat-engine
Restart=on-failure
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=30
WatchdogSec=120
StartLimitBurst=3
StartLimitIntervalSec=120

User=root
Group=root

# HARDENED: low priority — threat engine is observer, not actor.
Nice=10
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
CPUWeight=10

# Resource budget — threat engine must not starve normal workloads
CPUQuota=8%
MemoryMax=128M
MemoryHigh=96M
TasksMax=32

# State persistence
ReadWritePaths=/var/lib/hisnos /var/log/hisnos /run/hisnos

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
ReadOnlyPaths=/proc /sys/fs/cgroup
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_PACKET
RestrictNamespaces=yes
SystemCallFilter=@system-service @network-io @privileged
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_PTRACE CAP_AUDIT_READ
AmbientCapabilities=CAP_NET_ADMIN CAP_AUDIT_READ

# Journald integration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hisnos-threat

[Install]
# HARDENED: defer to graphical.target — must not block boot.
WantedBy=graphical.target
