# systemd/hisnos-stability-guard.timer
[Unit]
Description=HisnOS 6-hour Stability Guard Timer

[Timer]
OnBootSec=1h
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
